# articles/n8n-aws-cloudflare-selfhost/ec2/n8n-secrets.service
[Unit]
Description=Load n8n secrets from SSM into tmpfs
After=network-online.target
Wants=network-online.target
Before=docker-n8n.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/mkdir -p /run/n8n-secrets
ExecStartPre=/bin/mount -t tmpfs -o size=1M,mode=0700 tmpfs /run/n8n-secrets
ExecStart=/opt/n8n/load-secrets.sh
ExecStop=/bin/umount /run/n8n-secrets

[Install]
WantedBy=multi-user.target
